# node_filesystem/etc/service/powerstrip-calico/run
#!/bin/sh
cd /calico_containers/adapter/
exec python powerstrip_calico.py