# master/getting-started/mesos/vagrant-centos/units/docker.service.d/docker.conf
[Service]
EnvironmentFile=/etc/sysconfig/calico
ExecStart=
ExecStart=/usr/bin/dockerd --cluster-store=etcd://${ETCD_AUTHORITY}
