# node_fs/etc/service/powerstrip-calico/run
#!/bin/sh
cd /calico_docker/adapter/
exec python powerstrip_calico.py