# rpm/calico-felix.upstart
description "Calico Felix Agent"

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [S016]

respawn

exec su -s /bin/sh -c "exec /usr/bin/calico-felix --config-file /etc/calico/felix.cfg" root
