# calico_node/rpm/calico-dockerless-bird.service
[Unit]
Description=Bird for Dockerless Calico
After=calico-dockerless-confd.service
Requires=calico-dockerless-confd.service

[Service]
ExecStart=/usr/local/bin/bird -R -s /run/bird.ctl -d -c /etc/calico/confd/config/bird.cfg
EnvironmentFile=/etc/calico/calico-environment

[Install]
WantedBy=multi-user.target
