# node/etc/service/bird6/run
#!/bin/sh
exec bird6 -R -s bird6.ctl -d -c /config/bird6.cfg >>/var/log/calico/bird6.log 2>&1