# node_filesystem/etc/service/confd/run
#!/bin/sh
exec 2>&1
exec /confd -confdir=/ -interval=5 -watch --log-level=debug -node ${ETCD_AUTHORITY}