# node_filesystem/etc/service/confd/run
#!/bin/sh
exec /confd -confdir=/ -interval=5 -watch --log-level=debug -node ${ETCD_AUTHORITY} >>/var/log/calico/confd.log 2>&1