# node_filesystem/etc/service/calico-driver/run
#!/bin/sh
cd /calico_containers/
exec python docker_plugin.py