# node_filesystem/etc/service/calico-driver/run
#!/bin/sh
cd /calico_containers/driver/
exec python plugin.py