# v2.0/getting-started/mesos/vagrant/units/docker.service.d/docker.conf
[Service]
EnvironmentFile=/etc/sysconfig/calico
ExecStart=/usr/bin/dockerd --cluster-store=etcd://172.18.18.101:2379
