# node/root_overlay/etc/service/felix/run
#!/bin/sh
exec calico-felix --config-file=/config/felix.cfg