# debian/calico-acl-manager.upstart
# vim:set ft=upstart ts=2 et:
description "Calico ACL Manager"
author "Project Calico Maintainers <maintainers@projectcalico.org>"

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [!2345]

respawn

chdir /var/run

pre-start script
  mkdir -p /var/run/calico
  chown root:root /var/run/calico
end script

exec start-stop-daemon --start --chuid root --exec /usr/bin/calico-acl-manager -- --config-file=/etc/calico/acl_manager.cfg
