# v3.3/getting-started/mesos/installation/vagrant-centos/mesos-dns.service
[Unit]
Description=mesos-dns

[Service]
ExecStart=/usr/bin/mesos-dns -config=/etc/sysconfig/mesos-dns

[Install]
WantedBy=multi-user.target
