# node_filesystem/etc/rc.local
#!/bin/bash
# Record current boot time. Helps with debugging.
date > /tmp/boottime.txt

# Run Confd in onetime mode, to ensure that we have a working config in place to allow bird(s) and
# felix to start.  Don't fail startup if this confd execution fails.
./confd -confdir=. -onetime -node ${ETCD_AUTHORITY} || true