# node_filesystem/etc/service/powerstrip/run
#!/bin/sh
cd /powerstrip
exec twistd -noy powerstrip.tac