# docs/mesos/vagrant-centos/units/docker.service.d/docker.conf
[Service]
EnvironmentFile=/etc/sysconfig/calico
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// --cluster-store=etcd://${ETCD_AUTHORITY}